# etc/service/powerdns/run
#!/bin/bash -e

# PowerDNS sends everything to stderr, so we redirect.
exec /usr/sbin/pdns_server --disable-syslog=yes --log-timestamp=no --daemon=no 2>&1
