# etc/service/powerdns/run
#!/bin/bash -e

trap 'echo "ERROR: $BASH_SOURCE:$LINENO $BASH_COMMAND" >&2' ERR

# PowerDNS sends everything to stderr, so we redirect.
exec /usr/sbin/pdns_server --daemon=no 2>&1
