# etc-3.3/service/powerdns/run
#!/bin/bash -e

# PowerDNS sends everything to stderr, so we redirect.
exec /usr/sbin/pdns_server --daemon=no 2>&1
